# deploy/systemd/wt-agent@.service
[Unit]
Description=WovenTeam wt-agent (%i)
Documentation=file:///woventeam/repos/WovenTeam/docs/cli/wt-agent.md
After=network-online.target wt-roomd.service
Wants=network-online.target wt-roomd.service

[Service]
Type=exec
User=woven
Group=woven
WorkingDirectory=/woventeam/repos/WovenTeam
ExecStart=/woventeam/repos/WovenTeam/bin/wt-agent --agent %i --loop --config /woventeam/repos/WovenTeam/config/woventeam-phase0.conf
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wt-agent-%i

# Light hardening. The native room agent writes only under /woventeam: build
# artifacts when the launcher rebuilds, the JSONL room log, and state files.
NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/woventeam
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
